# jiractl/install/jiractl.conf
start on runlevel [2345]

script
  exec sudo -u jiractl ~jiractl/bin/node ~jiractl/apps/jiractl/jiractl.js >> /var/log/jiractl/jiractl.log 2>&1
end script

post-start script
  sleep 10
end script
